# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.


modprobe ip_gre
modprobe ipip

sh /etc/net.conf
sh /opt/jack2/servicios.conf
echo " Jack2 2.1 " > /etc/issue.net
cat /dev/null > /etc/motd
exit 0
